# gvisor-hello.service.tpl
[Service]
ExecStart=/usr/local/bin/runsc --debug --debug-log=PWD/gvisor-logs/ --fsgofer-host-uds run -bundle=PWD/bundle hello
Type=simple
StandardOutput=journal
